# systemd/ros2-chassis.service
[Unit]
Description=ROS2 Base Chassis Driver (no RViz)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zhengtuo
Group=zhengtuo
Environment=ROS_DOMAIN_ID=0
# 确保TTY可用（如需要）：添加 udev 规则后再启用
ExecStart=/home/zhengtuo/nailong_ros2_ws/scripts/start_chassis.sh
Restart=always
RestartSec=3
# 将日志输出到journalctl
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
